# crontab
# min   hour    day     month   weekday 	user 	command
*/5    *       *       *       *       root	/usr/bin/php /var/www/html/admin/cli/cron.php
